# snips-handler-weather.service
[Unit]
Description=Snips Handler WW
After=network.target

[Service]
User=respeaker
Group=respeaker
ExecStart=/usr/bin/python /usr/share/snips/handler/handler-weather.py 
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
